# rootfs/etc/services.d/smb/run
#!/usr/bin/with-contenv sh

samba --version

if [ ! -z "${DEBUG}" ]; then
    OPTS="${OPTS:-} --debuglevel=5"
fi

exec smbd --foreground --no-process-group --log-stdout ${OPTS:-}
